# scripts/netguard-api.service
[Unit]
Description=NetGuard Backend API
Documentation=https://github.com/YOUR_USERNAME/YOUR_REPO
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netguard-api
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=netguard-api
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
